# deploy/lightsail/dunpilot-api.service.example
[Unit]
Description=DunPilot API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/dunpilot/dnf-hell-optimizer-react
EnvironmentFile=/opt/dunpilot/dunpilot-api.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/dunpilot/dnf-hell-optimizer-react/.venv/bin/python neople_hell_api_server.py --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/dunpilot/dnf-hell-optimizer-react

[Install]
WantedBy=multi-user.target
